# cron/user_log_mailer.cron
# Run the user log mailer script on the 20th of every month at 11:59 PM


59 23 20 * * /usr/bin/python3 /home/aaron/user-log-mailer/scripts/generate_report.py >> /home/aaron/user-log-mailer/logs/cron_output.log

